# luci-app-service-weburl/root/etc/init.d/service_weburl
#!/bin/sh /etc/rc.common

START=99
USE_PROCD=1
PROG=/usr/bin/sqlite3
DB_PATH=/etc/service_weburl/data.db

init_database() {
    if [ ! -f "$DB_PATH" ]; then
        mkdir -p /etc/service_weburl
        $PROG "$DB_PATH" <<EOF
CREATE TABLE services (
    id INTEGER PRIMARY KEY AUTOINCREMENT,
    title TEXT NOT NULL,
    url TEXT NOT NULL,
    description TEXT,
    created_at TIMESTAMP DEFAULT CURRENT_TIMESTAMP,
    updated_at TIMESTAMP DEFAULT CURRENT_TIMESTAMP
);

CREATE TABLE logs (
    id INTEGER PRIMARY KEY AUTOINCREMENT,
    action TEXT NOT NULL,
    details TEXT,
    created_at TIMESTAMP DEFAULT CURRENT_TIMESTAMP
);
EOF
        chmod 644 "$DB_PATH"
    fi
}

start_service() {
    init_database
}

service_triggers() {
    procd_add_reload_trigger "service_weburl"
}